# layers/meta-appliance-os/recipes-core/rauc/files/appliance-mark-good.service
[Unit]
Description=Mark RAUC slot as good after successful boot
After=multi-user.target
ConditionPathExists=!/run/appliance-marked-good

[Service]
Type=oneshot
ExecStart=/usr/bin/rauc status mark-good
ExecStartPost=/usr/bin/touch /run/appliance-marked-good

[Install]
WantedBy=multi-user.target
