# layers/meta-appliance-apps/recipes-apps/appliance-reticulum/files/appliance-rnsd.service
[Unit]
Description=Reticulum transport daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStartPre=-/usr/bin/podman rm -f appliance-rnsd
ExecStart=/usr/bin/podman run \
    --name appliance-rnsd \
    --rm \
    --network host \
    -v /data/apps/reticulum/rnsd/config:/etc/rsReticulum:ro \
    -v /data/apps/reticulum/rnsd/storage:/var/lib/rsReticulum \
    --device /dev/ttyUSB0 \
    --device /dev/ttyACM0 \
    appliance-rnsd:latest
ExecStop=/usr/bin/podman stop -t 10 appliance-rnsd
Restart=always
RestartSec=5s

[Install]
WantedBy=multi-user.target
