# layers/meta-appliance-apps/recipes-apps/firefox-kiosk/files/home-kiosk-.mozilla-firefox-kiosk.mount
# Bind-mount persistent Firefox profile from /data into the kiosk user's
# home directory so bookmarks, cookies, and settings survive reboots
# and rootfs A/B updates.
#
# /data/apps/firefox-kiosk/profile/ is created by ExecStartPre in the app service.
# The mount point /home/kiosk/.mozilla/firefox/kiosk/ must exist on the rootfs.
#
# DefaultDependencies=no prevents the implicit Before=local-fs.target
# that would create a dependency cycle with systemd-tmpfiles-setup.

[Unit]
Description=Bind-mount Firefox Kiosk persistent profile
RequiresMountsFor=/data
DefaultDependencies=no
After=appliance-init.service data.mount
Before=firefox-kiosk.service

[Mount]
What=/data/apps/firefox-kiosk/profile
Where=/home/kiosk/.mozilla/firefox/kiosk
Type=none
Options=bind

[Install]
WantedBy=firefox-kiosk.service
